# packages/net-misc/dhcp/files/systemd/dhcpd@.service
# Copyright 2018- Julien Durillon <julien@durillon.xyz>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=DHCP Daemon on interface %I
Requires=network.target
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/dhcpd-%i.conf
ExecStart=/usr/host/bin/dhcpd -cf $DHCPD_CONF -lf /var/db/dhcpd-%i.leases -pf /run/dhcpd-%i.pid -q -f $DHCPD_OPTS %i

[Install]
WantedBy=multi-user.target
