# packages/net-misc/openntpd/files/openrc/initd/openntpd
#!/sbin/openrc-run

command="/usr/bin/openntpd"
pidfile="/run/openntpd.pid"
command_args="-p ${pidfile}"

depend() {
    need net
    use dns logger
}
