# packages/dev-lang/erlang/files/systemd/epmd@.socket
[Unit]
Description=Erlang Port Mapper Daemon Activation Socket

[Socket]
ListenStream=%I:4369
Accept=false

[Install]
WantedBy=sockets.target
