# packages/net-dns/bind/files/systemd/named.service
# Copyright 2010-2014 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Berkeley Internet Name Domain (DNS)
After=network.target

[Service]
ExecStart=/usr/bin/named -f -u named -n 1
ExecReload=/usr/bin/kill -HUP $MAINPID
PIDFile=/run/named/named.pid

PrivateTmp=true

[Install]
WantedBy=multi-user.target
