# packages/dev-db/postgresql/files/systemd/postgresql.service
[Unit]
Description=PostgreSQL @SLOT@ database server
Documentation=man:postgres(1)
After=network.target

[Service]
Type=notify
User=postgres
Environment=PGDATA=/var/lib/postgresql/@SLOT@/data
ExecStart=/usr/host/libexec/postgresql-@SLOT@/postgres -D ${PGDATA}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=0

[Install]
WantedBy=multi-user.target
